# cron/planefence.cron
*/2 * * * * root /home/pi/planefence/planefence.sh 2>&1
