# cron/planefence.cron
*/2 * * * * pi /usr/share/planefence/planefence.sh > /tmp/planefence.log 2>&1
