# cron/planefence.cron
*/2 * * * * pi /home/pi/planefence/planefence.sh > /tmp/planefence.log 2>&1
